# scripts/supabase-keepalive.service
[Unit]
Description=Supabase Keep-Alive Service
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/supabase-keepalive
ExecStart=/opt/supabase-keepalive/supabase-keepalive.sh
Restart=always
RestartSec=30
Environment=SUPABASE_ANON_KEY=your_actual_key_here

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=supabase-keepalive

[Install]
WantedBy=multi-user.target